# debian-pkg/nnext-server/DEBIAN/preinst
#!/bin/sh
# preinst script for nnext-server

# Copyright NNext Co. and/or licensed to NNext Co.

set -ex

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

stop_service_if_exists() {
    SERVICE=$1
    NUM=`set +e;invoke-rc.d --quiet --query $SERVICE stop;echo $?`
    if [ $NUM -eq 104 ]; then
        invoke-rc.d $SERVICE stop
    fi
}


case "$1" in
    install)
        if grep -E "^nnext:" /etc/passwd > /dev/null ; then
                echo "⚠️ nnext user already exist"
        else
                useradd -r -s /bin/sh -c "NNext Server" -M nnext
                echo "create nnext user"
        fi
        exit 0
    ;;
    upgrade)
        stop_service_if_exists nnext-server
        exit 0
    ;;
    abort-upgrade)
    ;;

    *)
        set +x
        echo "❗️ preinst called with unknown argument \`$1'" >&2
        set +x
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

